# test/wordpress/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command = sh -c "/usr/sbin/php-fpm7 -R -F && kill -s SIGTERM $(cat supervisord.pid)"

[program:nginx]
command = sh -c "nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"

# [program:mariadb]
# command = sh -c "/usr/bin/mysqld --datadir=/var/lib/mysql --socket=/run/mysqld/mysqld.sock && kill -s SIGTERM $(cat supervisord.pid)"

# [program:mariadb]
# command = sh -c "/usr/bin/mysqld_safe && kill -s SIGTERM $(cat supervisord.pid)"