# test/grafana/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command = sh -c "/var/www/html/telegraf-1.17.3/usr/bin/telegraf --config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"

[program:grafana]
command = sh -c "./grafana-server && kill -s SIGTERM $(cat supervisord.pid)"