# test/influxdb/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command = sh -c "./telegraf-1.17.3/usr/bin/telegraf --config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"

# в терминале ./grafana-7.3.7/bin/grafana-server НЕ ПОЛУЧАЛОСЬ
# [program:grafana]
# command = sh -c "./grafana-7.3.7/bin/grafana-server && kill -s SIGTERM $(cat supervisord.pid)"

[program:influxdb]
command = sh -c "influxd && kill -s SIGTERM $(cat supervisord.pid)"
